# data/user/.config/autostart/chromix-too-server.desktop
[Desktop Entry]
Type=Application
Name=Chromix Too Server
Exec=chromix-too-server
NoDisplay=true
